# supervisord_extra.ini
[program:dns_update]
command=/bin/python /the/path/to/oray_instead/tencent_dns_api.py --secret_id xxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxx --secret_key xxxxxxxxxxxxxxxxxxxxxxx --root_domain example.com --host transpond
directory=/the/path/to/oray_instead/
user=root
autorestart=true
redirect_stderr=true
