# supervisord_extra.ini
[program:dns_update]
command = python3 /the/path/to/oray_instead/ddns_sync.py
directory = /the/path/to/oray_instead/
user = root
autorestart = true
redirect_stderr = true
stdout_logfile = /the/path/to/oray_instead/log/ddns_sync_stdout.log
