# supervisord_extra.ini
[program:dns_update]
command = /bin/python /the/path/to/oray_instead/ddns_sync.py --secret_id xxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxx --secret_key xxxxxxxxxxxxxxxxxxxxxxx --root_domain example.com --host transpond --receivers "email1@example.com;email2@example.com"
directory = /the/path/to/oray_instead/
user = root
autorestart = true
redirect_stderr = true
stdout_logfile = /the/path/to/oray_instead/log/ddns_sync_stdout.log
